# config/includes.chroot/etc/init.d/anti-mitm
#!/sbin/openrc-run

# Определяем переменные
command="/usr/bin/python3"
command_args="/opt/anti-mitm.py"
pidfile="/run/anti-mitm.pid"
screen_session="anti-mitm"

depend() {
    # Укажите зависимости, если они есть
    after localmount
}

start() {
    ebegin "Starting anti-mitm service"
    
    # Проверяем, запущен ли уже процесс
    if pidof -o %PPID -x "$command $command_args" > /dev/null; then
        eend 0 "Service is already running"
        return 0
    fi

    # Запускаем в фоновом режиме с помощью screen
    /usr/bin/screen -dmS "$screen_session" "$command" "$command_args"
    
    # Записываем PID в файл
    echo $! > "$pidfile"
    
    eend $? "Service started"
}

stop() {
    ebegin "Stopping anti-mitm service"
    
    # Останавливаем сессию screen
    if screen -list | grep -q "$screen_session"; then
        screen -S "$screen_session" -X quit
        rm -f "$pidfile"  # Удаляем pid файл
        eend 0 "Service stopped"
    else
        eend 0 "Service is not running"
    fi
}
