# config/includes.chroot/etc/init.d/update-random-seed
#!/sbin/openrc-run
depend() {
    before shutdown
}

start() {
    ebegin "Updating random seed"
    dd bs=512 seek=34 count=1 oflag=direct if=/dev/urandom of=/dev/bilibop
    eend $?
}

stop() {
    # No stop action required
    return 0
}
   
